# conf/frappe-bench.target
[Unit]
After=network.target
Requires=__APP__-web.target __APP__-workers.target __APP__-redis.target

[Install]
WantedBy=multi-user.target
